# .docker/crontab
SHELL=/bin/bash

0 */4 * * * docker exec php-fpm bin/console app:job1
0 */12 * * * docker exec php-fpm bin/console app:job2
0 */24 * * * docker exec php-fpm bin/console app:job3